# RaspberryPi/yocto/src/meta-acm/recipes-core/udev-rules/files/acm-hci.rules
# try with a sleep 1 just before hciconfig

ACTION=="add", KERNEL=="hci0", RUN+="/usr/bin/hciconfig hci0 piscan"
